# packaging/targets/trusty/etc/cron.d/datim-auto-cert-updater
# check whether cert needs updated once a day, at 3:15 AM
15 3 * * * openhim_cert_updater /usr/share/datim-auto-cert-updater/check_for_renewal.sh
